# sysext/usr/lib/systemd/system/nvidia-mig-setup.service
[Unit]
Description=NVIDIA MIG Instance Setup
# Activated at boot by a TrueNAS PREINIT entry (midclt initshutdownscript),
# not by WantedBy=multi-user.target — sysext-shipped Install symlinks aren't
# reliably honored by TrueNAS's boot path. The script itself also waits up
# to 60s for /usr/bin/nvidia-smi in case the stock sysext hasn't merged yet.
After=systemd-sysext.service local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nvidia-mig-setup
