# sysext/usr/lib/systemd/system/nvidia-mig-setup.service
[Unit]
Description=NVIDIA MIG Instance Setup
# After= here is best-effort; on TrueNAS the actual sysext merge timing may
# differ. The script itself waits for /usr/bin/nvidia-smi (up to 60s) and
# exits cleanly with a log line if it never appears. We deliberately do NOT
# use ConditionPathExists, which would silently skip with no journal output.
After=systemd-sysext.service nvidia-persistenced.service local-fs.target
Before=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nvidia-mig-setup

[Install]
WantedBy=multi-user.target
